# to_home_dir/.config/systemd/user/gbrain-daily-resync.service.d/pool-size.conf
# Cap gbrain's postgres.js pool to avoid bursting the shared Supabase session
# pooler (default is 10 per process; multiple gbrain processes share one pooler).
# Shipped by brent-install. A drop-in (not an edit to the generated unit/script),
# so /setup-gbrain regeneration leaves it intact. Applies once the gbrain sync
# timers exist; harmless/inert before then.
[Service]
Environment=GBRAIN_POOL_SIZE=2
